# prometheus-borg-exporter.timer
[Unit]
Description=Prometheus Borg Exporter Timer

[Timer]
OnCalendar=*-*-* 08:00:00

[Install]
WantedBy=timers.target
